# leptographic.service
[Unit]
Description=Leptographic - Leptos UI System
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/leptographic
ExecStart=/var/www/leptographic/leptos-radix-ui
Environment=LEPTOS_SITE_ADDR=127.0.0.1:3000
Environment=LEPTOS_SITE_ROOT=/var/www/leptographic
Environment=RUST_LOG=info
Restart=always
RestartSec=5

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/leptographic

[Install]
WantedBy=multi-user.target
